# controllers/docker/controller/supervisord.conf
[supervisord]
nodaemon=true

[program:controller]
command=/usr/bin/controller
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
